# supervisor/mongo.conf
[program:mongo]
command=/opt/mongo/bin/mongod -f /opt/mongo/conf/mongo.conf
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
priority=5
