# supervisor/mysql.conf
[program:mysql]
command=/usr/sbin/mysqld --defaults-file=/etc/mysql/mysql.cnf
user=root
autostart=false
autorestart=true
startsecs=10
stopwaitsecs=30
priority=3
